# pyjamhole-supervisord.conf
[program:pyjamhole]
command=/usr/bin/python3 /path/to/pyjamhole.py